# etc/inittab
::sysinit:/etc/init.d/rcS
ttySAC0::respawn:/bin/sh

